# config/etc/services.d/fluentd/run
#!/usr/bin/with-contenv sh

fluentd -c /etc/fluentd/fluent.conf
